# deployment/systemd/samchat-registration-bot.service.d/ctt-shadow.conf.example
[Service]
Environment=CTT_RESPONSES_ROLLOUT=shadow
Environment=CTT_SHADOW_REVIEW_HANDOFF=on
